# net-analyzer/sinfo/files/sinfod.initd
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /usr/local/ssd/gentoo-x86/output/net-analyzer/cvs-repo/gentoo-x86/net-analyzer/sinfo/files/sinfod.initd,v 1.1 2011/11/28 22:20:22 radhermit Exp $

command="/usr/sbin/${SVCNAME}"
command_args="${SINFOD_OPTS}"

depend() {
	need net
}
